# src/systemd/bbfw.network.service.in
[Unit]
Description=Bindle Binaries Linux Network
Documentation=man:ip(8)
Wants=network.target
After=network-pre.target dbus.service
Before=network.target network.service

[Service]
Type=oneshot
ExecStart=@sbindir@/bbnet start
ExecStop=@sbindir@/bbnet stop
ExecReload=@sbindir@/bbnet restart
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

